# es.ini
[program:es]
command=su - elasticsearch --shell=/bin/bash -c "/usr/share/elasticsearch/bin/elasticsearch"
user=root
autorestart=true
stdout_logfile=/var/log/elasticsearch/elasticsearch.log
stderr_logfile=/var/log/elasticsearch/elasticsearch.log
