# es.ini
[program:es]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
autorestart=true
stdout_logfile=/var/log/elasticsearch/elasticsearch.log
stderr_logfile=/var/log/elasticsearch/elasticsearch.log
