# php-fpm.ini
[program:php-fpm]
command=bash -c "sleep 1 && php-fpm"  ; the program (relative uses PATH, can take args)
autostart=true                ; start at supervisord start (default: true)
autorestart=true        ; whether/when to restart (default: unexpected)
startretries=5                ; max # of serial start failures (default 3)
exitcodes=0,2,70                 ; 'expected' exit codes for process (default 0,2)
stopsignal=QUIT               ; signal used to kill process (default TERM)
stopwaitsecs=2                 ; max num secs to wait b4 SIGKILL (default 10)