# shadow.ini
[program:shadow]
stdout_logfile=/var/log/shadow/access.log
stderr_logfile=/var/log/shadow/error.log
autostart=true
autorestart=true