# shadow.ini
[program:shadow]
stdout_logfile=/var/log/shadow_access.log
stderr_logfile=/var/log/shadow_error.log
autostart=true
autorestart=true