# zookeeper/zookeeper.ini
[program:zookeeper]
command=/opt/zkServer.sh start-foreground
autorestart=true
stopsignal=KILL"
