# zkui/zkui.ini
[program:zkui]
command=/opt/zkui.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/zkui_err.log
stdout_logfile=/var/log/zkui.log
