# kafka/kfk.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh  /opt/kafka/config/server.properties
autostart=true
autorestart=true
stopwaitsecs=1
stopasgroup=true
startsecs=50
user=root
stderr_logfile=/var/log/kafkaerr.log
stdout_logfile=/var/log/kafka.log
